# kitchen/cookbooks/identity-analytics/templates/reporting-workers.target.erb
# Dropped off by chef
# Systemd target for reporting-workers
[Unit]
Description=reporting-workers
Requires=<%= @worker_service_descriptors %>

[Install]
WantedBy=multi-user.target
