# kitchen/cookbooks/login_dot_gov/templates/default/idp-workers.target.erb
# Dropped off by chef
# Systemd target for idp-workers
[Unit]
Description=idp-workers
Requires=<%= @worker_service_descriptors %>

[Install]
WantedBy=multi-user.target
